# debian/spine-publisher.upstart
# spine-publisher - Spine Configball Publisher
#

description	"Spine Configball Publisher"

start on filesystem
stop on runlevel [06]

expect fork
respawn

script

    SPINE_PUB_BIN="/usr/bin/spine-publisher"
    SPINE_PUB_OPTIONS="-c /etc/spine-mgmt/spine-publisher.conf"

    if [ -f /etc/default/spine-publisher ]; then
        . /etc/default/spine-publisher
    fi

    exec $SPINE_PUB_BIN $SPINE_BIN_OPTIONS
end script
